# daemons/opendoge-imu.service
[Unit]
Description=OpenDoge DM-IMU-L1 Bridge
Documentation=https://github.com/Lain-Ego/OpenDoge_firmware
After=default.target

[Service]
Type=simple

# 自动检测 DM-IMU-L1 (USB ID 6877:4d55) — 支持 ttyACM 和 ttyUSB
# 等待最多 15s 直到设备出现
ExecStartPre=/bin/sh -c ' \
for i in $$(seq 1 30); do \
    for dev in /dev/serial/by-id/usb-DM-Tech_DM-IMU-L1*; do \
        [ -e "$$dev" ] && echo "$$dev" && exit 0; \
    done; \
    sleep 0.5; \
done; \
echo "DM-IMU-L1 (6877:4d55) not found after 15s" >&2; \
exit 1'

# 启动 bridge
ExecStart=/bin/sh -c ' \
device=$$(ls /dev/serial/by-id/usb-DM-Tech_DM-IMU-L1* 2>/dev/null | head -1); \
if [ -z "$$device" ]; then \
    echo "DM-IMU-L1 device not found" >&2; \
    exit 1; \
fi; \
exec /usr/bin/python3 /home/orangepi/Opendoge/OpenDoge_firmware/hardware/imu/dm_imu_bridge.py \
    --source serial \
    --device "$$device" \
    --baud 921600 \
    --configure-usb \
    --output /tmp/opendoge_imu.state'

Restart=always
RestartSec=3
StartLimitIntervalSec=0
TimeoutStartSec=45

[Install]
WantedBy=default.target
