# scripts/opendoge-joystick.service
[Unit]
Description=OpenDoge Joystick Bridge
Documentation=https://github.com/Lain-Ego/OpenDoge_firmware
# User service cannot depend on system service directly — we wait for
# /dev/input/js0 in ExecStartPre instead.
After=default.target

[Service]
Type=simple
# Wait up to 30s for xboxdrv to create /dev/input/js0
ExecStartPre=/bin/sh -c 'for i in $$(seq 1 60); do [ -e /dev/input/js0 ] && exit 0; sleep 0.5; done; echo "js0 not found after 30s" >&2; exit 1'
ExecStart=/usr/bin/python3 /home/orangepi/Opendoge/OpenDoge_firmware/daemons/command_bridge/xbox_command_bridge.py \
    --output /tmp/opendoge_command.state \
    --require-rb
# Always restart — bridge has its own reconnection logic for js0 loss,
# but a hard crash (e.g. Python interpreter abort) still needs systemd.
Restart=always
RestartSec=3
StartLimitIntervalSec=0
TimeoutStartSec=45

[Install]
WantedBy=default.target
