# scripts/opendoge-joystick.service
[Unit]
Description=OpenDoge Joystick Bridge
Documentation=https://github.com/Lain-Ego/OpenDoge_firmware
After=default.target

[Service]
Type=simple
# Runs as the local user — xboxdrv must already provide /dev/input/js0.
# Use opendoge-xboxdrv.service (system-level, see install_services.sh) to
# auto-start xboxdrv on boot if needed.
ExecStart=/usr/bin/python3 /home/orangepi/Opendoge/OpenDoge_firmware/daemons/command_bridge/xbox_command_bridge.py \
    --output /tmp/opendoge_command.state \
    --require-rb
Restart=always
RestartSec=5
TimeoutStartSec=10

[Install]
WantedBy=default.target
